# ansible/roles/user/files/systemd/switch-gnome-wallpaper.timer
[Unit]
Description=Switch wallpaper each hours

[Timer]
OnCalendar=*:00:00
Unit=switch-gnome-wallpaper.service

[Install]
WantedBy=timers.target
